# deploy/xray-worker.service
[Unit]
Description=X-ray Data Collection Crawler/Processing Worker (Python)
After=network.target postgresql.service ollama.service
# ollama.service is optional; if it isn't installed the worker just uses the
# regex fallback. Listing it here only affects start ordering when present.
Wants=postgresql.service ollama.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/xray-data-collection/services/worker
ExecStart=/home/pi/xray-data-collection/services/worker/.venv/bin/python main.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
